# cpuctl.service
[Unit]
Description=cpuctl

[Service]
Type=idle
ExecStart=/usr/bin/cpuctl

[Install]
WantedBy = multi-user.target
